# kubernetes1.2/roles/haproxy/files/etc/init/cloudone_slb.conf
description "Cloudone slb daemon to control haproxy"
author "edcyc.tw@gmail.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /srv/cloudone_slb/cloudone_slb